# deploy/linux/junglelaw-server.service
[Unit]
Description=Jungle Law ENet dedicated server
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=junglelaw
Group=junglelaw
WorkingDirectory=/opt/junglelaw
StateDirectory=junglelaw
StateDirectoryMode=0750
Environment=HOME=/var/lib/junglelaw
Environment=XDG_DATA_HOME=/var/lib/junglelaw/.local/share
Environment=ZHANCHENG_BIND_HOST=0.0.0.0
Environment=ZHANCHENG_SERVER_PORT=24567
Environment=ZHANCHENG_MAX_CLIENTS=64
ExecStart=/opt/junglelaw/JungleLawServer.x86_64 --headless -- --bind-host=0.0.0.0 --server-port=24567 --max-clients=64
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/junglelaw
UMask=0077
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
